# blue_bench_eval/prompts/p1-09.yaml
id: p1-09
category: forensics
title: Wazuh Alert Contextualization
question: |
  As a SOC analyst, interpret the following sequence of Wazuh alerts from a single host over a 30-minute window. Determine whether this represents an active incident, a benign administrative event, or a false positive cluster. Justify your assessment.

  {"timestamp": "2026-01-15T14:01:00.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737284084.602", "data": {"srcip": "185.220.101.42", "srcport": "56129", "dstuser": "oracle"}}
  {"timestamp": "2026-01-15T14:01:08.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737158896.338", "data": {"srcip": "185.220.101.42", "srcport": "41243", "dstuser": "root"}}
  {"timestamp": "2026-01-15T14:01:16.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737120896.639", "data": {"srcip": "185.220.101.42", "srcport": "59715", "dstuser": "user"}}
  {"timestamp": "2026-01-15T14:01:24.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737447848.728", "data": {"srcip": "185.220.101.42", "srcport": "55662", "dstuser": "admin"}}
  {"timestamp": "2026-01-15T14:01:32.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737457311.834", "data": {"srcip": "185.220.101.42", "srcport": "46601", "dstuser": "test"}}
  {"timestamp": "2026-01-15T14:01:40.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737162448.122", "data": {"srcip": "185.220.101.42", "srcport": "45815", "dstuser": "postgres"}}
  {"timestamp": "2026-01-15T14:01:48.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737911982.213", "data": {"srcip": "185.220.101.42", "srcport": "44840", "dstuser": "test"}}
  {"timestamp": "2026-01-15T14:01:56.000+0000", "rule": {"id": "5710", "level": 5, "description": "sshd: Attempt to login using a non-existent user.", "groups": ["syslog", "sshd", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737176054.482", "data": {"srcip": "185.220.101.42", "srcport": "57245", "dstuser": "postgres"}}
  {"timestamp": "2026-01-15T14:02:09.000+0000", "rule": {"id": "5715", "level": 3, "description": "sshd: authentication success.", "groups": ["syslog", "sshd", "authentication_success"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737714483.203", "data": {"srcip": "185.220.101.42", "srcport": "52962", "dstuser": "deploy"}}
  {"timestamp": "2026-01-15T14:03:20.000+0000", "rule": {"id": "550", "level": 7, "description": "Integrity checksum changed for: '/etc/passwd'", "groups": ["ossec", "syscheck", "syscheck_file"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737452909.409", "data": {"file": "/etc/passwd", "changed_attributes": ["mtime", "md5", "sha1", "sha256"], "event": "modified"}}
  {"timestamp": "2026-01-15T14:03:35.000+0000", "rule": {"id": "550", "level": 7, "description": "Integrity checksum changed for: '/etc/shadow'", "groups": ["ossec", "syscheck", "syscheck_file"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737438273.239", "data": {"file": "/etc/shadow", "changed_attributes": ["mtime", "md5", "sha1", "sha256"], "event": "modified"}}
  {"timestamp": "2026-01-15T14:03:50.000+0000", "rule": {"id": "550", "level": 7, "description": "Integrity checksum changed for: '/etc/sudoers'", "groups": ["ossec", "syscheck", "syscheck_file"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737265154.849", "data": {"file": "/etc/sudoers", "changed_attributes": ["mtime", "md5", "sha1", "sha256"], "event": "modified"}}
  {"timestamp": "2026-01-15T14:05:00.000+0000", "rule": {"id": "5401", "level": 10, "description": "PAM: User login failed.", "groups": ["pam", "syslog", "authentication_failed"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737558089.910", "data": {"srcuser": "deploy", "dstuser": "root"}}
  {"timestamp": "2026-01-15T14:05:10.000+0000", "rule": {"id": "5402", "level": 3, "description": "PAM: Login session opened.", "groups": ["pam", "syslog", "authentication_success"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737770318.598", "data": {"srcuser": "deploy", "dstuser": "root"}}
  {"timestamp": "2026-01-15T14:05:20.000+0000", "rule": {"id": "5551", "level": 10, "description": "User successfully changed UID to root via sudo.", "groups": ["syslog", "sudo"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737780714.424", "data": {"srcuser": "deploy", "command": "/bin/bash"}}
  {"timestamp": "2026-01-15T14:07:30.000+0000", "rule": {"id": "2902", "level": 7, "description": "New systemd service created.", "groups": ["systemd", "service_creation"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737282664.824", "data": {"service_name": "system-health-check.service", "service_path": "/etc/systemd/system/system-health-check.service", "exec_start": "/tmp/.hidden/beacon"}}
  {"timestamp": "2026-01-15T14:07:40.000+0000", "rule": {"id": "510", "level": 7, "description": "Anomalous process execution detected.", "groups": ["ossec", "rootcheck"]}, "agent": {"id": "003", "name": "web-server-01", "ip": "10.10.5.22"}, "manager": {"name": "wazuh-manager"}, "id": "1737689275.810", "data": {"process": "/tmp/.hidden/beacon", "user": "root", "parent": "systemd"}}
expected_tools: []
expected_findings:
  - synonyms: ["brute force", "SSH brute force", "password spray", "credential stuffing", "185.220.101.42"]
  - synonyms: ["successful login", "authentication success", "deploy user", "5715"]
  - synonyms: ["/etc/passwd modified", "/etc/shadow modified", "/etc/sudoers modified", "credential files", "privilege escalation setup"]
  - synonyms: ["privilege escalation", "sudo to root", "UID changed to root", "5551", "escalation"]
  - synonyms: ["system-health-check.service", "persistence", "systemd service", "T1543.002"]
  - synonyms: ["/tmp/.hidden/beacon", "beacon process", "suspicious process", "malicious process in /tmp"]
  - synonyms: ["ACTIVE INCIDENT", "active incident", "confirmed incident", "not benign", "not false positive", "true positive"]
pass_criteria: Identifies the sequence as an active incident — brute force leading to successful login, credential file tampering, privilege escalation, persistence via systemd service, and beacon process execution
tags: [forensics, wazuh, incident-response, linux, phase1]
max_turns: 1
